# fancygotchi/mod/etc/systemd/system/pwny-restart.service
[Unit]
Description=Search syslog for Go map crash

[Service]
ExecStart=/root/pwny-restart.sh

[Install]
WantedBy=multi-user.target