# fancygotchi/mod/etc/systemd/system/pwny-restart2.service
[Unit]
Description=Search syslog for Go panic

[Service]
ExecStart=/root/pwny-restart2.sh

[Install]
WantedBy=multi-user.target